# conf/supervisor-app.ini
[program:app-uwsgi]
command = /usr/bin/uwsgi --ini /home/docker/code/conf/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
